# _mobile/ops/systemd/nanobot-gateway.service
# FR25 draft service: keep runtime auto-restart active and finalize
# user/path/command values during VPS deployment (p2-1, p4-3).

[Unit]
Description=Robotville Nanobot Gateway
Documentation=file:///opt/robotville/BMAD/_bmad/rbtv/_mobile/_docs/robotville-vps-access.md
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple

# Replace with the dedicated service account created in p2-1.
User=nanobot
Group=nanobot

# Replace with the absolute workspace path on the VPS.
WorkingDirectory=/opt/robotville/BMAD

# Nanobot reads ~/.nanobot/config.json; ensure HOME points to config location.
Environment=HOME=/srv/nanobot

# Secrets must stay outside git (NFR4). Create this file on VPS with chmod 600.
EnvironmentFile=/etc/robotville/nanobot-gateway.env

# Run nanobot gateway directly (ensure /usr/local/bin is in PATH for nanobot user)
ExecStart=/usr/local/bin/nanobot gateway

# FR25: runtime process must auto-restart on failure.
Restart=always
RestartSec=15

# Reliability guardrails.
TimeoutStartSec=120
TimeoutStopSec=45
KillMode=control-group

# Basic hardening compatible with Node-based runtime processes.
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=nanobot-gateway

[Install]
WantedBy=multi-user.target
