# ignite/deploy/rbtv-ignite.service
# rbtv-ignite.service — machine-agnostic systemd USER unit template for the ignite daemon.
#
# The daemon runs in place from the ignite module source tree; no deploy-root copy
# and no transport step.  Per-instance values are substituted at install time.
#
# Install-time substitution points (replace each @@NAME@@ with the instance value):
#   @@WORKSPACE_ROOT@@    - workspace root: the folder that roots .rbtv/
#   @@NODE_BIN@@          - absolute path to the Node.js 24 binary
#   @@IGNITE_SRC@@        - absolute path to the ignite module source root
#   @@WORKDIR_ROOT@@      - default workdir root every spawn falls back to (REQUIRED, D22:
#                           the seed config ships only a placeholder for this, so the daemon
#                           refuses to boot until a real per-instance value is supplied)
#   @@ENVIRONMENT_FILE@@  - optional EnvironmentFile path (use /dev/null to disable)
#
# Per-instance values reach the daemon through this unit and/or the EnvironmentFile —
# NEVER by editing config/spawn-profiles.yaml in the repo, which would write an instance
# path into the code tree (D26(3)).  EnvironmentFile= settings override Environment=
# (systemd.exec), so the file is the operator's override surface; RBTV_IGNITE_DATA_ROOT
# and RBTV_IGNITE_CARRIER are set there when the seed defaults do not fit.
#
# StateDirectory= provisions %S/rbtv-ignite (the runtime data root) for the user unit.
# In a user unit %S resolves to $XDG_STATE_HOME (~/.local/state).  The unit exports this
# path as RBTV_IGNITE_DATA_ROOT so the seed config's system-centric /var/lib/rbtv-ignite
# is overridden without editing the committed config.  An operator who moves data_root off
# the default via the EnvironmentFile owns creating that path.
#
# Install location for the user manager:
#   ~/.config/systemd/user/rbtv-ignite.service
#
# After installing and substituting:
#   systemctl --user daemon-reload
#   systemctl --user enable rbtv-ignite
#   systemctl --user start rbtv-ignite
#
# Do NOT commit a substituted version of this file.

# network-online.target does NOT exist for the user manager (it ships under
# /usr/lib/systemd/system, which the user manager never searches), so the After=/Wants=
# pair this unit carried as a SYSTEM unit resolved to `not-found` here and ordered nothing.
# It is dropped rather than left as a dependency that reads like a guarantee and is not one.
# The v1 bind is loopback (config bind.host 127.0.0.1), which needs no network wait. When
# p5-2 moves the bind to the tailnet, network readiness needs a REAL mechanism (bind retry
# in the daemon, or a user-manager-visible readiness unit) — not this line.
[Unit]
Description=RBTV ignite daemon

[Service]
Type=simple
WorkingDirectory=@@WORKSPACE_ROOT@@
EnvironmentFile=-@@ENVIRONMENT_FILE@@
Environment="RBTV_IGNITE_SRC=@@IGNITE_SRC@@"
Environment="RBTV_IGNITE_WORKSPACE_ROOT=@@WORKSPACE_ROOT@@"
Environment="RBTV_IGNITE_CONFIG_PATH=@@IGNITE_SRC@@/config/spawn-profiles.yaml"
Environment="RBTV_IGNITE_WORKDIR_ROOT=@@WORKDIR_ROOT@@"
Environment="RBTV_IGNITE_DATA_ROOT=%S/rbtv-ignite"
Environment="RBTV_IGNITE_USER_MANAGER=true"
StateDirectory=rbtv-ignite
ExecStart=@@NODE_BIN@@ @@IGNITE_SRC@@/server/index.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
